# ds4-disable-led/etc/udev/rules.d/60-ds4-disable-led.rules
#Turn off ds4 led

ACTION=="add", KERNEL=="js[0-9]*", SUBSYSTEM=="input", RUN+="/usr/bin/bash /usr/bin/ds4-disable-led"
